# config/monit.conf
# Monit starter - Added by Chetan

description "Monit - process monitoring daemon"

start on runlevel [2345]
exec /usr/bin/monit -v -Ic /etc/monit/monitrc
respawn
console output
